# 99-csm-device-remove.rules
SUBSYSTEM=="mhi", ACTION=="remove", KERNEL=="mhi[0-9]*_CSM_CTRL", RUN+="/usr/bin/csm-nbdkit-stop.sh %k"
