# csm-host-logger@.timer
[Unit]
Description=Host Logger Timer

[Timer]
OnActiveSec=5min
Unit=csm-host-logger@%i.service

[Install]
WantedBy=timers.target
