# csm-host-logger@.service
[Unit]
Description=Host log Collector

[Service]
Type=oneshot
ExecStart=/usr/bin/csm-host-logger.sh %i

[Install]
WantedBy=multi-user.target
